# FW2_3_0-HW_2_X_X_HW_3_X_X-sL/etc/init.d/avahi-daemon
#!/bin/sh /etc/rc.common
START=75

start() {
  avahi-daemon
}

stop() {
  killall avahi-daemon
}

reload() {
  killall -HUP avahi-daemon
}
