# 5500NAC_Firmware-V1.13.0/etc/init.d/lirc
#!/bin/sh /etc/rc.common
START=39

start() {
  mkdir -p /var/run/lirc
  lircd >&- 2>&- &
}

stop() {
  killall lircd >&- 2>&-
}
