# ReadyNASOS-6.10.5-x86_64/etc/etc/udev/rules.d/99-btrfs-scan.rules
SUBSYSTEM!="block", GOTO="btrfs_end"
KERNEL!="md*", KERNEL!="dm-*", GOTO="btrfs_end"

ACTION=="add|change",		ENV{ID_FS_TYPE}=="btrfs", \
	RUN+="/sbin/btrfs device scan $tempnode"

LABEL="btrfs_end"
