# ReadyNASOS-6.10.5-x86_64/etc/etc/udev/rules.d/99-ptp-camera.rules
#This is for auto media backup from cameras
ACTION=="add", ENV{ID_USB_INTERFACES}=="*:060101:*", \
	RUN+="/frontview/bin/fvbackup -s 1"
