# LSS5500NAC_Firmware-V1.13.0/etc/init.d/rclocal
#!/bin/sh /etc/rc.common

START=99

start() {
  [ -e /etc/rc.local ] && sh /etc/rc.local
}
