# 5500NAC_Firmware-V1.13.0/usr/lib/opkg/info/redis.list
/etc/redis.conf
/etc/init.d/memtest
/usr/sbin/redis-server
