# LSS5500NAC_Firmware-V1.13.0/usr/lib/opkg/info/avahi-mini.list
/usr/sbin/avahi-daemon
/etc/init.d/avahi-daemon
/etc/avahi/services/http.service
/usr/sbin/avahi-browse
