# LSS5500NAC_Firmware-V1.13.0/usr/lib/opkg/info/vsftpd.list
/etc/config/vsftpd
/lib/upgrade/keep.d/vsftpd
/etc/vsftpd.conf.dist
/usr/share/vsftpd/conf.lua
/usr/sbin/vsftpd
/etc/vsftpd.users
/etc/init.d/vsftpd
