# ReadyNASOS-6.10.5-arm/etc/etc/udev/rules.d/99-network-mtu.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="eth*" \
	RUN+="/lib/udev/set_mtu"
