# ReadyNASOS-6.10.5-arm/etc/etc/udev/rules.d/99-network-vlan.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="eth*", ENV{DEVTYPE}!="vlan" \
	RUN+="/lib/udev/create_vlans"
