# ReadyNASOS-6.10.5-arm/etc/etc/udev/rules.d/99-psu.rules
ACTION=="add", DRIVER=="pmbus", RUN+="/lib/udev/hotplug_event PSU"
