# 5500NAC2_Firmware_V1.13.0/usr/lib/opkg/info/genohm-scada.list
/www/scada/resources/images/clock.png
/usr/sbin/usbreset
/www/scada/modules/40-buildings.js.gz
/www/scada/resources/lib.js.gz
/lib/genohm-scada/web/vis/export.lua
/www/scada/resources/images/vis-order.png
/lib/genohm-scada/web/vis/usermode.lua
/www/scada/resources/images/user.png
/www/scada/resources/images/controls/dimmer.png
/www/scada/resources/images/controls/spinner.png
/usr/share/lua/genohm-scada/platform.lua
/lib/genohm-scada/web/scripting/print.lua
/lib/genohm-scada/web/scripting/save.lua
/lib/genohm-scada/cron/scheduler.lua
/lib/genohm-scada/web/visgraphics/css.lua
/lib/genohm-scada/web/objects/ws.lua
/etc/uci-defaults/01_cron_fixup
/lib/genohm-scada/cron/trends.lua
/lib/genohm-scada/web/remote/grp.lua
/lib/genohm-scada/web/objects/save.lua
/lib/genohm-scada/web/buildings/duplicate.lua
/lib/genohm-scada/core/platform-init.lua
/lib/genohm-scada/web/schedulers/delete.lua
/lib/genohm-scada/web/remote-handler.lua
/lib/genohm-scada/web/users/config-form.lua
/www/scada/resources/images/run.png
/lib/genohm-scada/web/general/config-save.lua
/www/scada/resources/languages/dutch.js.gz
/www/scada/resources/languages/german.js.gz
/www/scada/resources/minizip.js.gz
/lib/genohm-scada/config.default
/www/scada/resources/images/system.png
/lib/genohm-scada/web/buildings/floor-save.lua
/www/scada/resources/images/lock.png
/www/scada/resources/images/setvalue.png
/lib/genohm-scada/web/trends/form.lua
/usr/share/netsrv/logger.lua
/lib/genohm-scada/web/vistpl/touch.lua
/lib/genohm-scada/web/users/form.lua
/lib/genohm-scada/web/schedulers/events-delete.lua
/etc/procmon/redis-server
/lib/genohm-scada/web/main/editor.lua
/lib/genohm-scada/web/scenes/sequence-main.lua
/www/scada/resources/images/inactive.png
/lib/genohm-scada/web/timezones.lua
/lib/genohm-scada/web/visgraphics/save.lua
/lib/genohm-scada/web/objectlogs/csv.lua
/lib/genohm-scada/web/general/identify.lua
/www/scada/resources/languages/chinese.js.gz
/lib/genohm-scada/dbcheck.lua
/lib/genohm-scada/web/buildings/layouts.lua
/www/scada/resources/images/objects-reset.png
/www/scada/resources/images/upload.png
/lib/genohm-scada/web/trends/delete.lua
/www/scada/resources/core.js.gz
/lib/genohm-scada/web/buildings/sortorder.lua
/www/scada/resources/images/userlib-big.png
/www/scada/resources/images/list-view.png
/lib/genohm-scada/web/general/sortorder.lua
/www/scada/resources/languages/italian.js.gz
/lib/genohm-scada/web/scripting/form.lua
/lib/genohm-scada/web/schedulers/holidays-delete.lua
/www/scada/resources/images/config.png
/www/scada/modules/99-about.js.gz
/lib/genohm-scada/web/scenes/run.lua
/etc/uci-defaults/02_sysctl
/lib/genohm-scada/brand
/www/scada/resources/images/export.png
/lib/genohm-scada/web/objects/form.lua
/lib/genohm-scada/web/objects/setenums.lua
/www/scada/vis/trends.js.gz
/lib/genohm-scada/web/schedulers/holidays-save.lua
/etc/nginx/scada.locations
/lib/genohm-scada/web/general/config-form.lua
/lib/genohm-scada/web/remote/objects.lua
/lib/genohm-scada/web/visgraphics/font-delete.lua
/lib/genohm-scada/web/vistpl/index.lua
/www/scada/resources/images/vis-touch.png
/lib/genohm-scada/web/buildings/floor-form.lua
/usr/share/lua/rpc.lua
/lib/genohm-scada/web/failmessage.lua
/lib/genohm-scada/web/buildings/delete.lua
/lib/genohm-scada/web/scripting/restore.lua
/lib/genohm-scada/web/vis/checkpin.lua
/lib/genohm-scada/web/scenes/sequence-updown.lua
/lib/genohm-scada/web/vis/graph.lua
/lib/genohm-scada/web/users/delete.lua
/lib/genohm-scada/web/lib.lua
/lib/genohm-scada/web/objects/import-neighbours.lua
/lib/genohm-scada/web/vis/buildings.lua
/www/scada/modules/55-visgraphics.js.gz
/www/scada/resources/images/active.png
/lib/genohm-scada/web/trends/fetch.lua
/lib/genohm-scada/web/scenes/status.lua
/www/scada/resources/images/print.png
/www/scada/resources/images/userlib.png
/lib/genohm-scada/web/alerts/main.lua
/lib/genohm-scada/web/scenes/sequence-setvalue.lua
/www/scada/resources/languages/portuguese.js.gz
/lib/genohm-scada/web/scenes/sequence-save.lua
/www/scada/resources/languages/bulgarian.js.gz
/www/scada/resources/images/resident.png
/usr/share/lua/redis.lua
/www/scada/resources/images/restore-big.png
/usr/share/lua/vprint.lua
/lib/genohm-scada/web/remote/errors.lua
/www/scada/resources/images/tools.png
/lib/genohm-scada/web/buildings/building-save.lua
/lib/genohm-scada/web/general/backup.sh
/lib/genohm-scada/web/schedulers/updown.lua
/lib/genohm-scada/upgrade/storage.lua
/www/scada/resources/images/controls/circularslider.png
/usr/sbin/fontfix
/www/scada/vis/core.js.gz
/lib/genohm-scada/fw/flash.lua
/lib/genohm-scada/reset.sh
/www/scada/resources/images/help.png
/lib/genohm-scada/reference
/lib/genohm-scada/core/scripting-scheduled.lua
/lib/genohm-scada/web/vis/import.lua
/lib/genohm-scada/web/schedulers/holidays-form.lua
/www/scada/resources/images/bullet.gif
/lib/genohm-scada/cron/stats.lua
/lib/genohm-scada/web/vistpl/schedulers.lua
/lib/genohm-scada/web/schedulers/status.lua
/lib/genohm-scada/web/schedulers/events-main.lua
/www/scada/vis/busdecode.js.gz
/www/scada/resources/images/controls/stringinput.png
/lib/genohm-scada/web/general/factory-reset.lua
/www/scada/modules/75-logs.js.gz
/www/scada/resources/languages/czech.js.gz
/www/scada/resources/languages/greek.js.gz
/lib/genohm-scada/web/objects/setvalue.lua
/lib/genohm-scada/cron/backup.sh
/lib/genohm-scada/core/ipblocker.lua
/lib/genohm-scada/web/schedulers/events-status.lua
/www/scada/modules/60-scripting.js.gz
/www/scada/resources/images/controls/slider.png
/lib/genohm-scada/web/compat-handler.lua
/lib/genohm-scada/web/scripting/delete.lua
/lib/genohm-scada/web/users/logs.lua
/www/scada/resources/images/customvals.png
/lib/genohm-scada/upgrade/schedulers.lua
/lib/genohm-scada/web/schedulers/duplicate.lua
/lib/genohm-scada/web/trends/updown.lua
/www/scada/resources/images/note.png
/www/scada/resources/languages/polish.js.gz
/lib/genohm-scada/web/vis/floor.lua
/lib/genohm-scada/web/remote/stats.lua
/www/scada/resources/codemirror.css.gz
/www/scada/resources/images/esf-import.png
/www/scada/resources/images/vis-camera.png
/www/scada/resources/images/plan.png
/lib/genohm-scada/web/scenes/main.lua
/lib/genohm-scada/web/vistpl/pin.lua
/www/scada/resources/images/controls/charinput.png
/lib/genohm-scada/web/general/plugin.lua
/www/scada/resources/images/vis-object.png
/lib/genohm-scada/web/schedulers/main.lua
/www/scada/resources/languages/danish.js.gz
/lib/genohm-scada/web/scripting/editor.lua
/lib/genohm-scada/web/general/sync.lua
/lib/genohm-scada/web/scenes/sequence-form.lua
/lib/genohm-scada/core/scripting-event.lua
/lib/genohm-scada/web/visgraphics/js.lua
/usr/share/lua/genohm-scada/initscript.lua
/www/scada/resources/images/folder-alt.png
/lib/genohm-scada/cron/cleanup.lua
/lib/genohm-scada/upgrade/auth.lua
/lib/genohm-scada/web/buildings/building-form.lua
/lib/genohm-scada/web/scenes/sequence-add.lua
/www/scada/vis/theme.css.gz
/usr/share/lua/dbenv.lua
/lib/genohm-scada/web/visgraphics/delete.lua
/www/scada/modules/15-objectlogs.js.gz
/lib/genohm-scada/web/logs/main.lua
/www/scada/ext/ext.js.gz
/www/scada/resources/languages/russian.js.gz
/www/scada/resources/languages/turkish.js.gz
/www/scada/modules/39-scenes.js.gz
/usr/share/lua/rpclocal.lua
/www/scada/ext/images/button/btn-warn.gif
/lib/genohm-scada/web/buildings/floor-plan-save.lua
/www/scada/modules/70-alerts.js.gz
/www/scada/modules/50-visualization.js.gz
/lib/genohm-scada/web/buildings/updown.lua
/www/scada/resources/images/edit.png
/lib/upgrade/keep.d/genohm-scada
/lib/genohm-scada/web/general/bus-set-ip.lua
/lib/genohm-scada/web/scripting/backup.lua
/lib/genohm-scada/web/vis/order-save.lua
/www/scada/resources/images/vis.png
/lib/genohm-scada/web/scripting/value.lua
/lib/genohm-scada/web/general/neighbours.lua
/www/scada/resources/images/duplicate.png
/www/scada/resources/images/import.png
/lib/genohm-scada/web/visgraphics/font-save.lua
/usr/share/lua/localbus.lua
/www/scada/modules/80-errorlog.js.gz
/lib/genohm-scada/web/main/index.lua
/www/scada/ext/images/button/btn-alt.gif
/lib/genohm-scada/dbinit.lua
/www/scada/resources/language.js.gz
/www/scada/resources/scripting-editor.js.gz
/lib/genohm-scada/web/users/status.lua
/usr/share/nginx/auth.lua
/www/scada/resources/images/grid.png
/lib/genohm-scada/web/general/update.lua
/lib/genohm-scada/web/schedulers/events-save.lua
/lib/genohm-scada/web/vis/graphics.lua
/www/scada/resources/images/controls/toggler.png
/lib/genohm-scada/core/read-init.lua
/www/scada/vis/layout.css.gz
/www/scada/resources/images/controls/checkbox.png
/www/scada/vis/flotr2.js.gz
/lib/genohm-scada/cron/trends.sh
/usr/share/lua/trends.lua
/www/scada/resources/ext-scada.js.gz
/www/scada/resources/images/event.png
/www/scada/modules/61-users.js.gz
/www/scada/resources/images/clear.png
/www/scada/resources/images/checkers.png
/www/scada/vis/dygraph.js.gz
/www/scada/vis/schedulers.js.gz
/lib/genohm-scada/web/objects/check.lua
/lib/genohm-scada/web/scenes/save.lua
/lib/genohm-scada/web/general/backup.lua
/lib/genohm-scada/web/trends/main.lua
/lib/genohm-scada/web/users/main.lua
/www/scada/resources/languages/french.js.gz
/lib/genohm-scada/core/scripting-resident.lua
/etc/init.d/platform
/www/scada/resources/languages/hungarian.js.gz
/www/scada/modules/35-schedulers.js.gz
/lib/genohm-scada/web/schedulers/save.lua
/lib/genohm-scada/web/general/datetime-save.lua
/lib/genohm-scada/upgrade/trends.lua
/usr/share/nginx/conf.lua
/lib/genohm-scada/web/vis-handler.lua
/www/scada/resources/codemirror.js.gz
/lib/genohm-scada/web/remote/alerts.lua
/lib/genohm-scada/web/schedulers/vis-request.lua
/www/scada/resources/images/restore.png
/lib/genohm-scada/web/general/restore.lua
/etc/procmon/groupmonitor
/lib/genohm-scada/web/buildings/main.lua
/www/scada/vis/controls.js.gz
/www/scada/resources/images/vis-floorlink.png
/www/scada/resources/images/accept.png
/lib/genohm-scada/web/scripting/status.lua
/usr/share/lua/copas.lua
/lib/genohm-scada/web/scenes/duplicate.lua
/lib/genohm-scada/web/vis/buildings-lp.lua
/www/scada/resources/images/up.png
/etc/init.d/genohm-scada
/lib/genohm-scada/web/home-handler.lua
/lib/genohm-scada/web/schedulers/events-form.lua
/www/scada/resources/images/blockly.png
/lib/genohm-scada/web/objects/setvisparams.lua
/www/scada/resources/images/backup-big.png
/lib/genohm-scada/core/init.lua
/lib/genohm-scada/web/vistpl/trends.lua
/lib/genohm-scada/web/general/usage.lua
/lib/genohm-scada/cron/trends-backup.lua
/lib/genohm-scada/web/scripting/main.lua
/www/scada/resources/languages/swedish.js.gz
/www/scada/vis/gauge.js.gz
/www/scada/resources/images/initscript.png
/www/scada/resources/images/scheduled.png
/www/scada/resources/images/delete.png
/lib/genohm-scada/web/objects/main.lua
/lib/genohm-scada/web/errorlog/main.lua
/www/scada/vis/icons.woff
/lib/genohm-scada/web/scripting/export.lua
/lib/genohm-scada/web/vis/fonts.lua
/lib/genohm-scada/upgrade/common.lua
/www/scada/resources/images/unlock.png
/lib/genohm-scada/web/handler.lua
/lib/genohm-scada/web/objectlogs/main.lua
/www/scada/resources/images/down.png
/lib/genohm-scada/web/scenes/form.lua
/www/scada/resources/images/update.png
/www/scada/resources/images/editor.png
/usr/share/lua/genohm-scada/eibdgm.lua
/www/scada/resources/images/controls/slider-v.png
/www/scada/resources/images/editor-disabled.png
/lib/genohm-scada/web/schedulers/form.lua
/lib/genohm-scada/web/general/datetime-form.lua
/www/scada/modules/10-objects.js.gz
/lib/genohm-scada/web/objects/delete.lua
/usr/share/lua/genohm-scada/userlib.default.lua
/lib/genohm-scada/web/users/config-save.lua
/lib/genohm-scada/web/scripting/duplicate.lua
/www/scada/vis/jquery.js.gz
/lib/genohm-scada/dbinit.sh
/www/scada/resources/images/vis-usermode.png
/lib/genohm-scada/upgrade/pincode.lua
/lib/genohm-scada/web/trends/save.lua
/lib/genohm-scada/web/general/reset.lua
/lib/genohm-scada/web/users/save.lua
/www/scada/resources/images/add.png
/lib/genohm-scada/web/scripting/source.lua
/usr/share/lua/genohm-scada.lua
/lib/genohm-scada/web/objects/update.lua
/www/scada/resources/images/backup.png
/www/scada/resources/languages/spanish.js.gz
/www/scada/modules/38-trends.js.gz
/lib/genohm-scada/web/scripting/run.lua
/www/scada/modules/60-utilities.js.gz
/lib/genohm-scada/web/visgraphics/clear.lua
/lib/genohm-scada/web/scenes/delete.lua
/www/scada/resources/images/folder.png
/lib/genohm-scada/version
/lib/genohm-scada/web/scenes/sequence-delete.lua
/lib/genohm-scada/cron/backup.lua
/lib/genohm-scada/web/scenes/savelive.lua
/lib/genohm-scada/web/schedulers/holidays-main.lua
