# FW2_1_0-HW_1_X_X-sL/usr/lib/opkg/info/gpiod.list
/lib/restore/restore.sh
/etc/init.d/gpiod-ready
/usr/sbin/gpiod
/etc/init.d/gpiod
/lib/restore/defaults.sh
