# FW2_1_0-HW_1_X_X-sL/usr/lib/opkg/info/usb-modeswitch.list
/lib/usbmode/data.lua
/sbin/usbmode
/lib/usbmode/handler.lua
/etc/hotplug.d/usb
/etc/init.d/usbmode
