# ReadyNASOS-6.10.5-x86_64/etc/etc/default/config/lib/systemd/system/apps.mount
[Unit]
Description=Apps Directory
Before=readynasd.service
ConditionPathIsDirectory=/data/.apps

[Mount]
What=/data/.apps
Where=/apps
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
