# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/antivirus-sched-scan.timer
[Unit]
Description=Anti-Virus schedule scan timer

[Timer]
OnBootSec=17m
OnUnitActiveSec=1d
Unit=antivirus-sched-scan.service
