# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/apache2-safe.service
[Unit]
Description=The Apache HTTP Server (Safe Mode)
Conflicts=apache2.service

[Service]
Type=forking
ExecStart=/usr/sbin/apache2 -f /etc/frontview/apache/safe-mode.conf
