# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/auto-fw-update.service
[Unit]
Description=ReadyNAS Firmware Update
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/firmware_update -U
