# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/clamav-freshclam.timer
[Unit]
Description=Anti-Virus Definition Update Timer
BindsTo=clamav-daemon.service

[Timer]
OnBootSec=17m
OnUnitActiveSec=1d
Unit=clamav-freshclam.service
