# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/fvbackup-q.service
[Unit]
Description=ReadyNAS Backup Service
After=readynasd.service

[Service]
ExecStart=/frontview/bin/fvbackup-q
#Restart=on-abort

[Install]
WantedBy=multi-user.target
