# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/logtruncate.timer
[Unit]
Description=Log Truncate Timer

[Timer]
OnBootSec=21m
OnUnitActiveSec=1d
Unit=logtruncate.service

[Install]
WantedBy=multi-user.target
