# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/logtruncate.service
[Unit]
Description=Log Truncate Service
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/frontview/bin/logtruncate
