# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/lsyncd.service
[Unit]
Description=lsyncd

[Service]
ExecStart=/usr/bin/lsyncd /etc/lsyncd/lsyncd.conf.lua

[Install]
WantedBy=multi-user.target
