# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/mdadm.service
[Unit]
Description=MD monitoring service
Wants=mdcsrepair.service

[Service]
ExecStart=/sbin/mdadm -F -s -y -p /frontview/bin/md_event

[Install]
WantedBy=multi-user.target
