# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/mdcsrepair.service
[Unit]
Description=MD repair service

[Service]
ExecStart=/usr/sbin/mdcsrepaird -v

[Install]
WantedBy=multi-user.target
