# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/oauth2-vault.service
[Unit]
Description=OAuth2 Vault Service

[Service]
ExecStart=/usr/sbin/oauth2-vault
BusName=com.readynas.OAuth2Vault
