# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/pilgrim.service
[Unit]
Description=ReadyNAS LCD splasher
Before=readynas-misc.service
DefaultDependencies=no

[Service]
ExecStart=/usr/sbin/pilgrim

[Install]
WantedBy=sysinit.target
