# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/pvcache-prune.service
[Unit]
Description=Preview Cache Pruning Service
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/frontview/bin/pvcache-prune.sh

