# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/pvcache-prune.timer
[Unit]
Description=Preview Cache Prune Service Timer

[Timer]
OnBootSec=70m
OnUnitActiveSec=1d
Unit=pvcache-prune.service
