# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/radar.timer
[Unit]
Description=Radar Update Timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=4h
Unit=radar.service

[Install]
WantedBy=multi-user.target
