# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/radar.service
[Unit]
Description=Radar Update

[Service]
Type=simple
ExecStart=/opt/readycloud/bin/radar.sh
