# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/rc.local.service
[Unit]
Description=/etc/rc.local
ConditionPathExists=/etc/rc.local
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/etc/rc.local

[Install]
WantedBy=default.target
