# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/readytivo.service
[Unit]
Description=TiVo recording archiver
Requires=avahi-daemon.service

[Service]
ExecStart=/usr/sbin/readytivod -S

[Install]
WantedBy=multi-user.target
