# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/remote-access.service
[Unit]
Description=Remote Access for Tech Support

[Service]
ExecStart=/frontview/bin/remote_access
