# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/replisync.service
[Unit]
Description=ReadyDR Snapshot Replication
Wants=network-online.target
After=local-fs.target network-online.target

[Service]
Type=notify
ExecStart=/usr/sbin/replisync -d
Restart=on-abort

[Install]
WantedBy=multi-user.target
