# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/rn-update.timer
[Unit]
Description=Timer for ReadyNAS Update Service

[Timer]
OnBootSec=3m
OnUnitActiveSec=1d
Unit=rn-update.service

[Install]
WantedBy=multi-user.target
