# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/rn-update.service
[Unit]
Description=ReadyNAS Update
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/frontview/bin/rn-update.sh
IOSchedulingPriority=7
Nice=5
