# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/rsync-avahi.service
[Unit]
Description=Rsync Server Bonjour Advertisement
BindsTo=rsync.service

[Service]
ExecStart=/usr/bin/avahi-publish -f -s %H _rsync._tcp 873

