# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/ssh-avahi.service
[Unit]
Description=SSH Server Bonjour Advertisement
BindsTo=ssh.service

[Service]
ExecStart=/usr/bin/avahi-publish -f -s %H _ssh._tcp 22
