# ReadyNASOS-6.10.5-arm/lib/lib/systemd/system/wasabi-s3.service
[Unit]
Description=Wasabi Web Service client for ReadyNAS
After=readynasd.service

[Service]
ExecStart=/usr/sbin/ng-csd -s Wasabi -l 1
Restart=on-failure

[Install]
WantedBy=multi-user.target
