# ReadyNASOS-6.10.5-arm/opt/replication/etc/cron.d/spool
0 */4 * * * root /bin/bash /opt/replication/bin/spoolp.sh
