# ReadyNASOS-6.10.5-arm/opt/replication/etc/cron.d/watchdog
*/1 * * * * root /bin/bash /opt/replication/etc/init.d/watchdog.sh
