# FW2_3_0-HW_2_X_X_HW_3_X_X-sL/usr/lib/opkg/info/bacnet.list
/usr/lib/libbacnet.so
/lib/bacnet/common.lua
/usr/sbin/bacwi
/usr/sbin/bacwp
/usr/sbin/bacepics
/lib/bacnet/handler.lua
/usr/sbin/bacnet
/etc/init.d/bacnet
/etc/config/bacnet
/usr/sbin/bacrp
