# FW2_3_0-HW_2_X_X_HW_3_X_X-sL/usr/lib/opkg/info/base-files.list
/etc/resolv.conf
/var
/etc/hotplug2-init.rules
/etc/preinit
/etc/init.d/network
/lib/upgrade/platform.sh
/etc/hotplug2-common.rules
/etc/shells
/etc/protocols
/sbin/ifdown
/etc/config/system
/usr/share/udhcpc/default.script
/etc/sysctl.conf
/sbin/sysupgrade
/etc/shadow
/etc/rc.common
/etc/init.d/watchdog
/etc/init.d/sysntpd
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/services
/lib/upgrade/common.sh
/lib/upgrade/keep.d/base-files
/etc/group
/etc/hotplug2.rules
/etc/hosts
/lib/functions.sh
/etc/config/network
/etc/init.d/rcS
/sbin/ifup
/etc/inittab
/etc/profile
/etc/init.d/mount
/etc/TZ
/lib/functions/service.sh
/lib/upgrade/keep.d/base-files-essential
/etc/passwd
/etc/mtab
/etc/init.d/boot
