# FW2_3_0-HW_2_X_X_HW_3_X_X-sL/usr/lib/opkg/info/eibd.list
/usr/bin/eibd
/etc/init.d/eibd
/lib/upgrade/keep.d/eibd
/etc/config/eibd
/usr/lib/eibd-probe.lua
/etc/procmon/eibd
